# deluge/build/conf/supervisor.conf
[supervisord]
logfile = /config/supervisord.log
logfile_maxbytes = 10MB
logfile_backups = 5
loglevel = debug
user = root
nodaemon = true
umask = 000

[include]
files = /etc/supervisor/conf.d/*.conf

[program:start-script]
autorestart = false
startsecs = 0
user = root
command = /root/start.sh
umask = 000

[program:watchdog-script]
autorestart = false
startsecs = 0
user = user
command = /home/user/start_user.sh
umask = 000